# study-log-api/scripts/com.studylog.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.studylog.api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/obst2580/github_project/personal/study-log/study-log-api/dist/index.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>STUDYLOG_PORT</key>
        <string>3100</string>
        <key>STUDYLOG_DB_PATH</key>
        <string>/Users/obst2580/.studylog/studylog.db</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/obst2580/.studylog/logs/api.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/obst2580/.studylog/logs/api-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/obst2580/github_project/personal/study-log/study-log-api</string>
</dict>
</plist>
